# systemd/trailcam-web.service
[Unit]
Description=Trailcam web interface (webapp.py)
After=network.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/home/YOUR_USERNAME/trailcam
ExecStart=/usr/bin/python3 /home/YOUR_USERNAME/trailcam/webapp.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
